# ipsw-patch/bundles/BootNeuter.bundle/files/System/Library/LaunchDaemons/com.apple.SpringBoard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.SpringBoard</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.springboard.migserver</key>
		<true/>
	</dict>
	<key>OnDemand</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/CoreServices/SpringBoard.app/SpringBoard</string>
	</array>
	<key>ServiceIPC</key>
	<false/>
	<key>UserName</key>
	<string>mobile</string>
</dict>
</plist>
